# .config/hypr/provisioning/hermes/vision/hermes-vision.service.in
[Unit]
Description=hermes-vision — screen-memory feed (gemma vision → honcho)
# gemma.service is the local llama.cpp server on :11434; honcho.service is the
# podman-compose memory stack on :8000. Neither is created by install.sh, so on
# a fresh machine both are absent — which is fine and deliberate:
#   After=  on a unit that does not exist is a no-op.
#   Wants=  on a unit that does not exist logs "not found" and carries on.
# Where they DO exist, this orders vision behind the model and the memory store
# instead of racing them at login. The daemon tolerates both being down anyway:
# HERMES_VISION_LLM / HONCHO_URL are just HTTP endpoints, a failed call is
# retried, and Restart=on-failure covers a hard exit.
After=graphical-session.target gemma.service honcho.service
Wants=gemma.service

[Service]
# Batch mode: buffer keyframes in memory and send them to gemma as ONE
# multi-image request per window, instead of one call per frame. The vision
# encoder amortises across images — measured at 1280px on an RX 9060 XT:
# 1 frame 1.9s, 2 frames 3.5s, 4 frames 4.0s, 8 frames 4.8s. Eight batched
# frames therefore cost ~4.8s against ~15.2s taken one at a time.
#
# Frames live in RAM only; the PNG grim writes is deleted as soon as it is
# encoded, so nothing screen-derived is ever persisted to disk.
#
# Set HERMES_VISION_BATCH_WINDOW=0 to restore one gemma call per frame.
Environment=HERMES_VISION_BATCH_WINDOW=900
Environment=HERMES_VISION_BATCH_MAX=8
Type=simple
# grim needs a live WAYLAND_DISPLAY. systemd's activation-time environment
# snapshot can predate uwsm importing it, so read it from the live user manager
# at exec time and wait for the socket before starting the daemon.
ExecStart=/bin/sh -c 'for i in $(seq 1 60); do WD="$(systemctl --user show-environment | sed -n "s/^WAYLAND_DISPLAY=//p")"; if [ -n "$WD" ] && [ -S "$XDG_RUNTIME_DIR/$WD" ]; then break; fi; sleep 1; done; export WAYLAND_DISPLAY="$WD"; echo "hermes-vision: WAYLAND_DISPLAY=$WD"; exec python3 %h/.hermes/bin/hermes-vision.py'
Restart=on-failure
RestartSec=10
# The daemon itself is light (grim + base64 + http); the heavy work is gemma on
# the GPU. A modest nice + a CPU ceiling keep even a bad loop from competing
# with interactive work.
Nice=5
CPUQuota=200%

[Install]
WantedBy=graphical-session.target
